# docs/systemd/timelapse-tapo-highres.service
[Unit]
Description=Timelapse high-resolution network camera
After=network.target

[Service]
Type=simple
User=username
WorkingDirectory=/home/username/repos/timelapse
ExecStart=/usr/bin/python3 /home/username/repos/timelapse/timelapse.py tapo.highres sunflowers-highres --flush 100
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
